# docker/systemd/didin-whatsapp-reconnect.service
[Unit]
Description=TikTrend Finder WhatsApp Reconnect Worker
Documentation=https://github.com/arkheioncorp/tiktrend-facil
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=didin
Group=didin
WorkingDirectory=/var/www/tiktrend-facil/backend

# Ambiente
Environment=PYTHONUNBUFFERED=1
Environment=ENVIRONMENT=production
Environment=RECONNECT_CHECK_INTERVAL=60
EnvironmentFile=-/etc/tiktrend-facil/whatsapp.env

# Comando
ExecStart=/var/www/tiktrend-facil/venv/bin/python -m workers.whatsapp_reconnect

# Restart automático
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Timeouts
TimeoutStartSec=30
TimeoutStopSec=15

# Recursos
MemoryMax=200M
CPUQuota=20%

# Segurança
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tiktrend-whatsapp-reconnect

[Install]
WantedBy=multi-user.target
